# systemd/restic-exporter.socket
[Unit]
Description=Prometheus Restic Exporter

[Socket]
ListenStream=9183

[Install]
WantedBy=sockets.target
